# deploy/systemd/kb-consolidate.timer
[Unit]
Description=Run KnotKB consolidation nightly at 03:30

[Timer]
OnCalendar=*-*-* 03:30:00
Persistent=true

[Install]
WantedBy=timers.target
