# deploy/systemd/kb-consolidate.service.tmpl
[Unit]
Description=KnotKB nightly consolidation (kb consolidate --apply)
Wants=docker.service

[Service]
Type=oneshot
ExecStart=@KB_RUNNER@
