# deploy/cron/kb-consolidate.crontab.tmpl
# kb-consolidate BEGIN (managed by deploy/install-scheduler.sh — do not edit)
30 3 * * * @KB_RUNNER@ >> @KB_LOG@.log 2>&1
# kb-consolidate END
